# networking/vlan-sync.service
[Unit]
Description=Sync VLANs from config.yaml
After=network.target
Wants=network.target

[Service]
Type=oneshot
# Set REPO_DIR to the repository path on the target device. Adjust if different.
Environment=REPO_DIR=/home/dietpi/tag-tapper-pi
Environment=VLAN_CONFIG=${REPO_DIR}/config.yaml
WorkingDirectory=${REPO_DIR}
ExecStart=${REPO_DIR}/networking/sync_vlans.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
